# srcs/requirements/click/conf/supervisor.conf
[supervisord]
nodaemon=true

[program:website]
command=/usr/bin/npm start
user=nobody
autostart=true
autorestart=true
stdout_logfile=/var/log/click.log
stderr_logfile=/var/log/click.err

